# rootfs/etc/profile.d/vlab.sh
if test -e /opt/vlab/path.sh; then
  . /opt/vlab/path.sh
fi
